# cron/crontab
#crontab
0 7 * * 1-5 root *notifications/max_limit_reached
